# sg.service
[Unit]
Description=Suckless graphs
After=network.target

[Service]
ExecStart=/var/sg/bin/sg-service

[Install]
WantedBy=multi-user.target
