# sgd@.service
[Unit]
Description=Suckless graphs

[Service]
ExecStart=/usr/local/bin/sgd
User=%i

[Install]
WantedBy=multi-user.target
