# examples/temp@.timer
[Unit]
Description=Get temperature hourly

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
